# j5a-nightshift/j5a-nightshift.service
[Unit]
Description=J5A Nightshift Queue Processor
Documentation=file:///home/johnny5/Johny5Alive/j5a-nightshift/README.md
After=network.target
StartLimitIntervalSec=3600
StartLimitBurst=3

[Service]
Type=oneshot
User=johnny5
Group=johnny5
WorkingDirectory=/home/johnny5/Johny5Alive/j5a-nightshift
Environment="HOME=/home/johnny5"

# Pre-flight safety checks
ExecStartPre=/home/johnny5/Johny5Alive/j5a-nightshift/ops/pre_flight_check.sh

# Dispatch inbox tasks to processing queue
ExecStartPre=/home/johnny5/Johny5Alive/scripts/nightshift_dispatcher.sh

# Main processing
ExecStart=/usr/bin/python3 /home/johnny5/Johny5Alive/j5a-nightshift/process_nightshift_queue.py

# Post-processing summary
ExecStartPost=/usr/bin/python3 /home/johnny5/Johny5Alive/j5a-nightshift/ops/nightshift_summary.py

# Logging
StandardOutput=append:/home/johnny5/Johny5Alive/j5a-nightshift/ops/logs/systemd.log
StandardError=append:/home/johnny5/Johny5Alive/j5a-nightshift/ops/logs/systemd_error.log

# Resource limits
MemoryMax=12G
CPUQuota=75%

# Restart on failure (but don't loop infinitely)
Restart=on-failure
RestartSec=300

[Install]
WantedBy=multi-user.target
